# tools/edge-chat.service
[Unit]
Description=JC1 Edge AI Chat UI
After=network-online.target edge-gateway.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/lucineer/.openclaw/workspace/tools/edge-chat.py --port 8081
Restart=on-failure
RestartSec=5
WorkingDirectory=/home/lucineer/.openclaw/workspace/tools
Environment=PATH=/home/lucineer/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=/home/lucineer

[Install]
WantedBy=default.target
